# deployment/systemd/nts-backup.service
[Unit]
Description=NTS daily SQLite backup
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/trading/news-trading-system
ExecStart=/home/trading/news-trading-system/.venv/bin/python3 /home/trading/news-trading-system/scripts/backup_db.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
